# darkhttpd.service
[Unit]
Description=Basic HTTP Daemon

[Service]
#User=p_client
ExecStart=/opt/client_repo/darkhttpd.sh

[Install]
WantedBy=multi-user.target 
